# supervisord.conf
[supervisord]
nodaemon=true

[program:acehttp]
command=/home/aceproxy/aceproxy/acehttp.py
user = aceproxy
directory = /home/aceproxy/aceproxy
autorestart = true
environment = HOME="/home/aceproxy"
